# files/profile.d/010-docs.sh
#!/usr/bin/env bash


PATH="${PATH}:\${HOME}/src/frontend/.scalingo/node/bin:\${HOME}/src/frontend/.scalingo/yarn/bin:\${HOME}/src/frontend/node_modules"
export PATH
